# src/nginx.cron
0 0 1 * * * root /usr/local/bin/aws-nginx-remoteip-update
